# system/etc/profile.d/homebin_override.sh
 export PATH=$HOME/.local/bin_override:$PATH
